# supervisord.conf
[program:memcached]
user = daemon
command=/usr/bin/memcached -m $MEMORY
