# supervisord.conf
[program:memcached]
command=/usr/bin/memcached -m $MEMORY
