# supervisor.conf
[program:memcached]
command=/usr/bin/memcached -m $MAX_MEMORY
